# config/messenger-worker.ini
[program:messenger-consume]
command=php /Users/weaverryan/Sites/messenger/bin/console messenger:consume async_priority_high async --time-limit=3600
user=weaverryan
numprocs=2
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
